# config/prometheus_exporter.conf

[program:prometheus_exporter]
command=/var/www/html/php-fpm_exporter server